# cron.d/alma-exports
# Cron jobs for alma exports
# Possible future cron jobs 
3 3 * * * root /mnt/alma/alma-scripts/scripts/Update-Export.sh 
30 0 2 * * root /mnt/alma/alma-scripts/scripts/Full-Export.sh