# cron.d/credit-card-slips
# Cron jobs for credit card slips
#
# Daily at 0801
01 08 * * * gituser /home/gituser/alma-scripts/scripts/Credit-Card-Slips.sh
#
#make sure to keep required newline at end of file -

