# cron.d/timdex-exports
# Cron jobs for TIMDEX exports

CURRENT_DATE=date +%Y%m

# Daily at 0853
53 08 * * * gituser cd /mnt/alma/alma-scripts && PATH=/usr/local/bin/pipenv:$PATH && pipenv run llama concat-timdex-export --export_type UPDATE >> /mnt/alma/logs/timdex-concat.log 2>&1

# At 0030 on the 2nd of the month (runs on full export from the 1st of the
# month)
30 0 2 * * gituser cd /mnt/alma/alma-scripts && PATH=/usr/local/bin/pipenv:$PATH && pipenv run llama concat-timdex-export --export_type FULL --date "$(${CURRENT_DATE})01" >> /mnt/alma/logs/timdex-concat.log 2>&1
