# cron.d/timdex-exports
# Cron jobs for TIMDEX exports
#
# Daily at 0853
53 08 * * * gituser /mnt/alma/alma-scripts/scripts/Update-Export.sh
#
# # At 0030 on the 2nd of the month (runs on full export from the 1st of the
# # month)
30 0 2 * * gituser /mnt/alma/alma-scripts/scripts/Full-Export.sh